# containers/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:agent-server]
command=/usr/local/bin/openhands-agent-server --port %(ENV_AGENT_SERVER_PORT)s
directory=/workspace
environment=SESSION_API_KEY="%(ENV_SESSION_API_KEY)s",OPENHANDS_SUPPRESS_BANNER="1"
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:automation]
command=uvicorn automation.app:app --host 0.0.0.0 --port 8000
directory=/app
environment=OPENHANDS_SUPPRESS_BANNER="1"
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
