# containers/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

[program:agent-server]
command=/agent-server/.venv/bin/python -m openhands.agent_server --port %(ENV_AGENT_SERVER_PORT)s
directory=/agent-server
environment=SESSION_API_KEY="%(ENV_SESSION_API_KEY)s"
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:automation]
command=uvicorn automation.app:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
