# containers/supervisord.combined.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ---- Nginx (main entry point on port 8000) ----
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5

# ---- OpenHands App Server (port 3000) ----
# Serves the GUI frontend and manages conversation sandboxes
[program:openhands]
command=/app/.venv/bin/python -m uvicorn openhands.server.listen:app --host 127.0.0.1 --port 3000
directory=/app
environment=PYTHONPATH="/app",VIRTUAL_ENV="/app/.venv",PATH="/app/.venv/bin:%(ENV_PATH)s",OPENHANDS_SUPPRESS_BANNER="1",RUNTIME="local",OH_PERSISTENCE_DIR="/data/openhands",SERVE_FRONTEND="false",SKIP_DEPENDENCY_CHECK="1"
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# ---- Agent Server for Automations (port 3002) ----
# This is the embedded agent server used by automation runs
# (separate from ProcessSandboxService which spawns its own)
[program:agent-server]
command=/usr/local/bin/openhands-agent-server --port %(ENV_AGENT_SERVER_PORT)s
directory=/workspace
environment=SESSION_API_KEY="%(ENV_SESSION_API_KEY)s",OPENHANDS_SUPPRESS_BANNER="1"
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=15

# ---- Automation Service (port 8001) ----
# Uses its own venv at /app/automation-venv to avoid SDK version conflicts
[program:automation]
command=/app/automation-venv/bin/uvicorn automation.app:app --host 127.0.0.1 --port 8001
directory=/app/automation-src
environment=PYTHONPATH="/app/automation-src",VIRTUAL_ENV="/app/automation-venv",PATH="/app/automation-venv/bin:%(ENV_PATH)s",OPENHANDS_SUPPRESS_BANNER="1"
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
