# containers/supervisord.combined.conf
# Supervisord configuration for combined Agent Server GUI + Automation service
#
# Services:
#   - nginx (port 8000) - reverse proxy
#   - agent-server (port 3002) - conversations + automation runs
#   - automation (port 8001) - automation API

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ---- Nginx (main entry point on port 8000) ----
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5

# ---- Agent Server (port 3002) ----
# Shared by both the GUI and automation runs
# Conversations are saved to ~/.openhands/conversations (mount your ~/.openhands)
[program:agent-server]
command=/usr/local/bin/openhands-agent-server --host 127.0.0.1 --port %(ENV_AGENT_SERVER_PORT)s
directory=/workspace
environment=OH_CONVERSATIONS_PATH="/root/.openhands/conversations",OPENHANDS_SUPPRESS_BANNER="1"
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# ---- Automation Service (port 8001) ----
[program:automation]
command=/app/automation-venv/bin/uvicorn automation.app:app --host 127.0.0.1 --port 8001
directory=/app/automation-src
environment=PYTHONPATH="/app/automation-src",VIRTUAL_ENV="/app/automation-venv",PATH="/app/automation-venv/bin:%(ENV_PATH)s",OPENHANDS_SUPPRESS_BANNER="1"
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=15
